# meta-cube/recipes-support/overc-conftools/source/overc-conftools.service
[Unit]
Description=OverC puppet system configuration
After=local-fs.target openvswitch-nonetwork.service network.target
Before=lxc.service

[Service]
Type=forking
Environment=FACTER_ipaddress="192.168.0.2"
ExecStart=/usr/bin/puppet apply /etc/puppet/manifests/site.pp

[Install]
WantedBy=multi-user.target
